# ansible/build-server/roles/rtc_wakeup/files/ fetch-holidays.service
[Unit]
Description=Fetch Japanese Holidays

[Service]
ExecStart=/usr/local/bin/fetch_holidays.sh
