# ansible/build-server/roles/rtc_wakeup/files/set-rtc-wakeup.timer
[Unit]
Description=Set RTC Wake Timer

[Timer]
OnCalendar=*-*-* 01:59:00
Persistent=true

[Install]
WantedBy=timers.target
