# ansible/build-server/roles/rtc_wakeup/files/set-rtc-wakeup.service
[Unit]
Description=Set RTC Wake Time
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/m-set-rtc-wakeup.sh
